# silktalk-pro/dist-package/scripts/silktalk.service
[Unit]
Description=SilkTalk Pro P2P Node
After=network.target

[Service]
Type=simple
User=silktalk
Group=silktalk
WorkingDirectory=/var/lib/silktalk
ExecStart=/usr/bin/silktalk start --port 4001 --ws-port 8080
Restart=always
RestartSec=10
Environment="SILKTALK_LOG_LEVEL=info"
Environment="SILKTALK_CONFIG_PATH=/etc/silktalk/config.json"
Environment="SILKTALK_DATA_PATH=/var/lib/silktalk"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/silktalk

[Install]
WantedBy=multi-user.target
